# script/run_services/payment_bridge.service
[Unit]
Description=Payment Bridge (MCP Backend) Service

[Service]
Type=simple
ExecStart=/home/zzq/go/src/payment-bridge/build/payment-bridge
StandardOutput=/home/zzq/go/src/payment-bridge/build/payment-bridge/build/pb_log.log
Restart=always
RestartSec=1
KillMode=mixed

[Install]
WantedBy=multi-user.target